# systemd/crt-autostart.service
[Unit]
Description=crt-tv boot rotation (weather 5 min, then the videos channel)
After=weather-kiosk.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tv autostart

[Install]
WantedBy=multi-user.target
