# systemd/crt-splash.service
[Unit]
Description=CRT-TV teletext boot splash
DefaultDependencies=no
After=local-fs.target
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/crt-tv/splash.sh

[Install]
WantedBy=sysinit.target
