# systemd/crt-splash.service
[Unit]
Description=CRT-TV teletext boot splash
DefaultDependencies=no
# console-setup: draw only after the final console font is loaded — a later
# setfont would remap the mosaic block glyphs already on screen
After=local-fs.target console-setup.service
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/crt-tv/splash.sh

[Install]
WantedBy=sysinit.target
