# systemd/ws3kp.service
[Unit]
Description=WeatherStar 3000+ local server
After=network-online.target
Wants=network-online.target

[Service]
User=crt
WorkingDirectory=/opt/ws3kp
Environment=WS3KP_PORT=8083
ExecStart=/usr/bin/node index.mjs
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
