# deploy/crt-tv.service
# crt-tv control service (FastAPI/uvicorn).
# Template — deploy/install.sh substitutes __CRT_TV_USER__ / __CRT_TV_DIR__ for
# the installing user and repo path, then installs it to /etc/systemd/system/.

[Unit]
Description=crt-tv control service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__CRT_TV_USER__
WorkingDirectory=__CRT_TV_DIR__
ExecStart=__CRT_TV_DIR__/.venv/bin/uvicorn crt_tv.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3
Environment=CRT_TV_CONFIG=__CRT_TV_DIR__/config.toml

[Install]
WantedBy=multi-user.target
